# wireguard/rootfs/etc/services.d/api/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community :WireGuard
# The most simple HTTP API you've ever seen.
# Provides status of WireGuard peers.
# ==============================================================================

                
